# config/DAVShare.conf
[program:DAVShare]
command=python /var/www/DAVShare/fsmon.py /var/www/webdav/
autorestart=true
user=root